# deploy/mctrader-collector.service
[Unit]
Description=mctrader forward-only OHLCV/tick/orderbook collector
Documentation=https://github.com/mclayer/mctrader-data
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mctrader
Group=mctrader

# State directory: /var/lib/mctrader/data (created by systemd if absent)
StateDirectory=mctrader/data
StateDirectoryMode=0750

Environment="MCTRADER_DATA_ROOT=/var/lib/mctrader/data"
Environment="PYTHONUNBUFFERED=1"

# Adjust path to your venv install (uv default: /opt/mctrader/.venv/bin/mctrader-data)
ExecStart=/opt/mctrader/.venv/bin/mctrader-data collect --top-n 10 --include transactions,orderbook --log-level INFO

# Graceful shutdown — SIGTERM + 30s drain before SIGKILL
KillSignal=SIGTERM
TimeoutStopSec=30s

# Restart policy
Restart=on-failure
RestartSec=10s
StartLimitIntervalSec=300
StartLimitBurst=5

# Resource hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ReadWritePaths=/var/lib/mctrader/data

# Network access required (Bithumb WebSocket)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Memory hint (top-10 collector ~ 100-200MB steady state)
MemoryMax=512M
TasksMax=64

[Install]
WantedBy=multi-user.target
